# digital-signage-saas/deploy/gunicorn.service.example
# Copy to /etc/systemd/system/signage.service
# sudo systemctl daemon-reload && sudo systemctl enable --now signage
#
# Create /etc/signage.env (chmod 600) with at least:
#   SECRET_KEY=...
#   DATABASE_URL=postgresql://user:pass@127.0.0.1/signage
#   DATA_DIR=/var/lib/signage-data
# Optional: BASE_PATH=, GOOGLE_*, STRIPE_*, MAIL_*

[Unit]
Description=Digital Signage SaaS (Gunicorn)
After=network.target postgresql.service

[Service]
User=www-data
Group=www-data
WorkingDirectory=/var/www/digital-signage-saas
EnvironmentFile=/etc/signage.env
ExecStart=/var/www/digital-signage-saas/venv/bin/gunicorn \
  --workers 4 \
  --worker-class gthread \
  --threads 2 \
  --bind 127.0.0.1:8000 \
  --timeout 90 \
  --graceful-timeout 30 \
  --keep-alive 5 \
  --max-requests 2000 \
  --max-requests-jitter 200 \
  app:app
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
